# containers/php8-fpm/devops/docker/supervisord/conf.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
# user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


; STARTUP

[program:composer-scripts]
command=/bin/sh -c "%(ENV_APP_DIR)s/bin/console doctrine:migrations:migrate -n && composer run-script post-install-cmd"
startsecs=0
autorestart=false
startretries=0
user=www-data

[program:tail-log]
command=/bin/sh -c "touch %(ENV_APP_DIR)s/var/log/prod.log && tail -f %(ENV_APP_DIR)s/var/log/prod.log"
user=www-data
numprocs=1
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=false
startretries=0
redirect_stderr=true
